# ansible/roles/chrome/templates/chromedriver.service.j2
[Service]
ExecStartPre=/usr/bin/env sleep 30s
ExecStart=/usr/bin/env chromedriver --port={{ chromedriver_port }}
Restart=always
