# ansible/roles/chrome/templates/chromedriver.service.j2
[Service]
User={{ ansible_user }}
ExecStart=chromedriver --port "{{ chromedriver_port }}"
Restart=always
RestartSec=30s
