# ansible/roles/ansible_pull_timer/templates/ansible-pull.timer.j2
[Timer]
Unit=ansible-pull.service
OnUnitActiveSec=5m

[Install]
WantedBy=default.target
